# swap-cleanup.service
[Unit]
Description=Swap Cleanup Service
Documentation=man:swap_cleanup.sh
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/swap-cleanup-tool/swap_cleanup.sh
User=root

[Install]
WantedBy=multi-user.target